# scripts/com.ccquota.usage-push.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS LaunchAgent 範本:每 5 分鐘把本機 7d/5h 用量推回 ccquota。
  安裝:
    1. 改下面的 __SCRIPT__ / __SERVER__ / __ACCOUNT__ / __INGEST_TOKEN__
    2. cp 到 ~/Library/LaunchAgents/com.ccquota.usage-push.plist
    3. chmod 600 ~/Library/LaunchAgents/com.ccquota.usage-push.plist  (內含 ingest token)
    4. launchctl load ~/Library/LaunchAgents/com.ccquota.usage-push.plist
  注意:要跑在「有真 Claude Code 登入此帳號」的機器上(它讀 Keychain 現成 token)。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ccquota.usage-push</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__SCRIPT__</string>
    <string>--server</string>
    <string>__SERVER__</string>
    <string>--account</string>
    <string>__ACCOUNT__</string>
    <string>--token</string>
    <string>__INGEST_TOKEN__</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/ccquota-usage-push.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ccquota-usage-push.log</string>
</dict>
</plist>
